# csg2@.service.in
[Unit]
Description=ChlorideSiteGenerator 2 - for site %i

[Service]
ExecStart=@POSTINSTALLPATH@/bin/csg2 /srv/csg2/%I

[Install]
WantedBy=multi-user.target
